# Web-Service/supervisord-single.conf
# Production Supervisord Configuration for Single Container Deployment
# Manages both Redis and HandPose application processes with:
# - Graceful shutdown handling (30s timeout)
# - Log rotation (50MB files, 10 backups)
# - Process monitoring and auto-restart
# - Structured logging with timestamps

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
loglevel=info
minfds=4096
minprocs=200

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# ============================================================================
# Redis Service (Cache & Session Store)
# ============================================================================
[program:redis]
command=/usr/bin/redis-server --bind 127.0.0.1 --port 6379 --maxmemory 256mb --maxmemory-policy allkeys-lru --appendonly yes --dir /app/data --save 300 10 --save 60 10000 --loglevel notice --timeout 60
user=nodejs
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/supervisor/redis-stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisor/redis-stderr.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
priority=10
startsecs=5
stopwaitsecs=30
stopsignal=TERM
stopasgroup=true
killasgroup=true

# ============================================================================
# HandPose Application (Web Server)
# ============================================================================
[program:handpose]
command=node /app/dist/index.js
directory=/app
user=nodejs
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/supervisor/handpose-stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisor/handpose-stderr.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
priority=20
startsecs=10
stopwaitsecs=30
stopsignal=TERM
stopasgroup=true
killasgroup=true
# Start after Redis is ready
depends_on=redis

environment=
    NODE_ENV="production",
    PORT="5000",
    HOST="0.0.0.0",
    DATABASE_URL="%(ENV_DATABASE_URL)s",
    JWT_SECRET="%(ENV_JWT_SECRET)s",
    JWT_EXPIRES_IN="%(ENV_JWT_EXPIRES_IN)s",
    REDIS_HOST="127.0.0.1",
    REDIS_PORT="6379",
    REDIS_DB="0",
    GCS_BUCKET_NAME="%(ENV_GCS_BUCKET_NAME)s",
    GCS_PROJECT_ID="%(ENV_GCS_PROJECT_ID)s",
    GOOGLE_APPLICATION_CREDENTIALS="%(ENV_GOOGLE_APPLICATION_CREDENTIALS)s",
    ENABLE_MOCK_STORAGE="%(ENV_ENABLE_MOCK_STORAGE)s",
    CORS_ORIGIN="%(ENV_CORS_ORIGIN)s",
    MAX_FILE_SIZE="%(ENV_MAX_FILE_SIZE)s",
    CLEANUP_ENABLED="%(ENV_CLEANUP_ENABLED)s",
    CLEANUP_RETENTION_DAYS="%(ENV_CLEANUP_RETENTION_DAYS)s",
    CLEANUP_CRON="%(ENV_CLEANUP_CRON)s",
    SMTP_HOST="%(ENV_SMTP_HOST)s",
    SMTP_PORT="%(ENV_SMTP_PORT)s",
    SMTP_SECURE="%(ENV_SMTP_SECURE)s",
    SMTP_USER="%(ENV_SMTP_USER)s",
    SMTP_PASSWORD="%(ENV_SMTP_PASSWORD)s",
    SMTP_FROM_SYSTEM="%(ENV_SMTP_FROM_SYSTEM)s",
    SMTP_FROM_NOREPLY="%(ENV_SMTP_FROM_NOREPLY)s",
    SMTP_FROM_ADMIN="%(ENV_SMTP_FROM_ADMIN)s",
    APP_URL="%(ENV_APP_URL)s",
    SUPABASE_URL="%(ENV_SUPABASE_URL)s",
    SUPABASE_ANON_KEY="%(ENV_SUPABASE_ANON_KEY)s",
    SUPABASE_SERVICE_ROLE_KEY="%(ENV_SUPABASE_SERVICE_ROLE_KEY)s"

# ============================================================================
# Automated Backup Scheduler (Runs every 6 hours)
# ============================================================================
[program:backup-scheduler]
command=/bin/sh /app/scripts/backup-scheduler.sh
user=nodejs
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/backup-scheduler-stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisor/backup-scheduler-stderr.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
priority=30
startsecs=0
stopwaitsecs=10
